# crontab
*/1 * * * * root echo hypriot/rpi-cron-example >>/var/log/cron.log 2>&1
